# deploy/systemd/trae-web.service
[Unit]
Description=Trae Web Service
After=network.target

[Service]
Type=simple
User=debian
Group=debian
WorkingDirectory=/home/debian/quantum-matrix
Environment=NODE_ENV=production
ExecStart=/usr/bin/npm run preview -- --host 0.0.0.0 --port 8080
Restart=on-failure
RestartSec=5
TimeoutStopSec=10

[Install]
WantedBy=multi-user.target
